# backup/etc/profile.d/xdg-base-dirs.sh
export \
	XDG_DATA_HOME="/home/default/.local/share" \
	XDG_CONFIG_HOME="/home/default/.config" \
	XDG_CACHE_HOME="/tmp/.cache/$LOGNAME" \
	XDG_STATE_HOME="/home/default/.local/state"
