# backup/etc/local.d/tap.start
#!/usr/bin/sh

/usr/bin/ip tuntap add mode tap tap0
/usr/bin/ip addr add 1.0.0.1/8 dev tap0
/usr/bin/ip link set dev tap0 up
